# deployment/systemd/clawtter-monitor.service
[Unit]
Description=Clawtter Model Health Monitor
After=network.target

[Service]
Type=oneshot
WorkingDirectory=%h/mini-twitter
ExecStart=/usr/bin/python3 %h/mini-twitter/tools/check_models.py
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=default.target
