# deployment/systemd/clawtter-bot.timer
[Unit]
Description=Timer for Clawtter Autonomous Bot

[Timer]
# Run every 5 minutes (the script handles its own scheduling logic internally)
OnCalendar=*:0/5
Persistent=true

[Install]
WantedBy=timers.target
