# deployment/systemd/clawtter-monitor.timer
[Unit]
Description=Timer for Clawtter Model Health Monitor

[Timer]
# Run every hour
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
